# board/nextthing/chip/overlay/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
null::sysinit:/bin/stty -F /dev/ttyS1 115200 &
null::sysinit:/bin/mount -t proc proc /proc
null::sysinit:/bin/mount -o remount,rw / &
null::sysinit:/bin/mkdir -p /dev/pts &
null::sysinit:/bin/mkdir -p /dev/shm &
null::sysinit:/bin/mount -a &
null::sysinit:/bin/hostname -F /etc/hostname &
null::sysinit:/bin/echo /sbin/mdev > /proc/sys/kernel/hotplug
null::sysinit:/sbin/mdev -s &
null::sysinit:/sbin/modprobe sun4i-codec
null::sysinit:/usr/bin/amixer -c 0 cset numid=6 4 &
null::sysinit:/opt/bin/blink-leds &
null::sysinit:/opt/bin/fade-pwms 0 0.5 &
null::sysinit:/opt/bin/fade-pwms 1 1 &
null::sysinit:/usr/sbin/alsactl restore
null::sysinit:/usr/bin/speaker-test -c 2 -t sine &

# now run any rc scripts
# ::sysinit:/etc/init.d/rcS

ttyS0::respawn:/bin/sh
ttyS1::respawn:/usr/bin/arecord -V "stereo" -r 8000 -c 2 -d 10 /dev/null

# Put a getty on the serial port
# ttyS1::respawn:/sbin/getty -L ttyS1 115200 xterm # GENERIC_SERIAL
# ttyGS0::respawn:/sbin/getty -L ttyGS0 115200 vt100 # gadget serial
# tty0::respawn:/sbin/getty -L tty0 115200 vt100 # framebuffer console

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/etc/init.d/rcK
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
